# etc/bg3po.crontab
# m h  dom mon dow   command
# auto reply to mail to bg3po
*/5 * * * * /home/bg3po/bg3po-scripts/bin/mail-autoreply.py > /home/bg3po/bg3po-scripts/log/mail-autoreply.log 2>//home/bg3po/bg3po-scripts/log/mail-autoreply_error.log

# auto post Monthly Bazaar
0 10 1 * * /home/bg3po/bg3po-scripts/bin/bazaar-autopost.py > /home/bg3po/bg3po-scripts/log/bazaar-auto-post.log 2>&1

# auto post Top 10 list
# look in /usr/local/bg3po/logs for logs.
# 0 7 15 * * /usr/local/bg3po/bin/top10-autopost.sh

# Repost old game of the weeks
0 10 * * * /usr/bin/python3 /home/bg3po/bg3po-scripts/bin/regotw.py > /home/bg3po/bg3po-scripts/log/regotw.log 2>&1

# post GOTW based on data in wiki/game_of_the_week
# every wednesday at 10 am localtime.
0 10 * * wed /usr/bin/python3 /home/bg3po/bg3po-scripts/bin/update_gotw.py -l debug > /home/bg3po/bg3po-scripts/log/gotw.log 2>&1

# testing.
*/1 * * * * echo cron job run at $(date) > /home/bg3po/bg3po-scripts/logs/test.log
